# supervisord.conf
[unix_http_server]
file=/djangodocker/djangodocker/daemons/supervisor.sock # path to your socket file

[supervisord]102
logfile=/djangodocker/logs/supervisord.log # supervisord log file
#logfile_maxbytes=50MB; maximum size of logfile before rotation
logfile_backups=10 # number of backed up logfiles
loglevel=debug # info, debug, warn, trace
pidfile=/djangodocker/djangodocker/daemons/supervisord.pid # pidfile location
nodaemon=false # run supervisord as a daemon
#minfds=1024; number of startup file descriptors
#minprocs=200; number of process descriptors
user=root # default user
childlogdir=/djangodocker/logs/ # where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

#[supervisorctl]
#serverurl=unix:///djangodocker/daemons/supervisor.sock; use unix:// schem for a unix sockets.

[include]

# Uncomment this line for celeryd for Python
files=djangod.conf

[program:djangoproject]
directory=/djangodocker
command=python3 manage.py runserver 0.0.0.0:8000
autostart=true
autorestart=true
priority=20